# crontab
*/30 * * * * /your/pkgcache/path/here/clear-indices.sh
